# recipes-bsp/swap-root/files/swap-root.service.in
[Unit]
Description=Verify successful boot after root swapped
Requires=multi-user.target
After=multi-user.target

[Service]
Type=oneshot
ExecCondition=/bin/sh -c "/usr/sbin/swap-root state | grep -q SWAPPING"
ExecStart=/usr/sbin/swap-root commit
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target