# recipes-bsp/nvram/nvram-check/nvram-check.service.in
[Unit]
Description=Check if nvram is online
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=@BINDIR@/nvram-check
RemainAfterExit=Yes

[Install]
WantedBy=multi-user.target
