# recipes-bsp/nvram/nvram-check/nvram-check.service.in
[Unit]
Description=Check if nvram is online
Requires=local-fs.target
After=local-fs.target
StartLimitIntervalSec=0

[Service]
StandardOutput=null
Type=oneshot
ExecStart=@BINDIR@/nvram
RemainAfterExit=Yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
