# recipes-bsp/nvram/nvram-hostname/nvram-hostname.service.in
[Unit]
Description=Set hostname from serial number in nvram
Requires=systemd-hostnamed.service nvram-check.service
After=systemd-hostnamed.service nvram-check.service

[Service]
Type=oneshot
ExecStart=/usr/bin/nvram-hostname @NVRAM_HOSTNAME@ @NVRAM_HOSTNAME_VARIABLE@
ExecStartPost=/bin/systemctl disable nvram-hostname.service
RemainAfterExit=Yes

[Install]
WantedBy=multi-user.target
