# recipes-bsp/nvram/files/nvram-factory.mount.in
[Unit]
Description=Production VPD partition
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target
ConditionPathExists=/dev/mmcblk0gp0p1

[Mount]
What=/dev/mmcblk0gp0p1
Where=/nvram/factory
Type=ext4
Options=@OPTS@

[Install]
WantedBy=local-fs.target
